# infra/lxc/scripts/config-manager/config-manager.service
[Unit]
Description=LXC Configuration Manager - Git-based config sync
Documentation=https://github.com/kethalia/infrahaus
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'printf "System is being configured by config-manager.\\nPlease wait for initial setup to complete.\\nCheck progress: journalctl -u config-manager -f\\n" > /var/run/nologin'
ExecStart=/usr/local/bin/config-sync.sh
ExecStopPost=/bin/rm -f /var/run/nologin
RemainAfterExit=no
Environment=HOME=/root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=config-manager

# Security considerations:
# This is a system configuration service that needs to:
# - Create/modify users (useradd/usermod)
# - Write to /etc (sudoers, systemd units, config files)
# - Install packages (Docker, etc. which may need to set capabilities)
# - Modify /usr (install binaries)
#
# Applying ProtectSystem or NoNewPrivileges breaks these operations.
# This service runs once during initial container setup, not continuously.
# Trade-off: Broad permissions for flexibility vs. security hardening.
RuntimeDirectory=config-manager

# Timeouts
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
