# infra/lxc/scripts/config-manager/config-manager.service
[Unit]
Description=LXC Configuration Manager - Git-based config sync
Documentation=https://github.com/kethalia/pve-home-lab
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/config-sync.sh
RemainAfterExit=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=config-manager

# Security hardening
# Note: ProtectSystem=full (not strict) because the file deployment engine
# needs to write configuration files to arbitrary system paths (e.g. /etc,
# /home). 'full' still protects /usr and /boot.
ProtectSystem=full
ReadWritePaths=/var/log/config-manager /opt/config-manager /etc/config-manager /var/lib/config-manager
RuntimeDirectory=config-manager
PrivateTmp=yes
NoNewPrivileges=yes

# Timeouts
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
